# conf/runit/412node/run
#!/bin/bash

exec 2>&1
exec sudo -u foronetwotrip node /home/foronetwotrip/412trip/index.js